# app-misc/systembus-notify/files/systembus-notify.service
[Service]
ExecStart=/usr/bin/systembus-notify

[Install]
WantedBy=default.target
